# ops/com.murmweb.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.murmweb.watchdog</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/mariuskerkvliet/money-with-claude-factory/ops/watchdog.sh</string>
  </array>

  <!-- Every 15 minutes: clears stuck locks and alerts on staleness. -->
  <key>StartInterval</key>
  <integer>900</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>WorkingDirectory</key>
  <string>/Users/mariuskerkvliet/money-with-claude-factory</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/mariuskerkvliet/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>ProcessType</key>
  <string>Background</string>

  <key>StandardOutPath</key>
  <string>/Users/mariuskerkvliet/money-with-claude-factory/ops/launchd.watchdog.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/mariuskerkvliet/money-with-claude-factory/ops/launchd.watchdog.err.log</string>
</dict>
</plist>
